# launchd/com.google.driveviewerexpiry.daily.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.google.driveviewerexpiry.daily</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>__PROJECT_ROOT__/scripts/run-daily-audit.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>~/Library/Logs/drive-viewer-expiry/run.log</string>
	<key>StandardErrorPath</key>
	<string>~/Library/Logs/drive-viewer-expiry/run.err</string>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<!--
	Optional: set folder + sheet for scheduled runs (uncomment and edit):
	<key>EnvironmentVariables</key>
	<dict>
		<key>FOLDER_ID</key>
		<string>your_drive_folder_id</string>
		<key>SPREADSHEET_ID</key>
		<string>your_sheet_id</string>
		<key>EXPIRY_HOURS</key>
		<string>24</string>
		<key>ACTIVITY_HOURS</key>
		<string>24</string>
	</dict>
	-->
</dict>
</plist>
